# scripts/example.feed-worker.cron
*/10 * * * * /path/to/morning-serial/scripts/feed-worker.js >> /path/to/logfile
